# templates/rtorrent/rsync.cron.tmpl
#!/bin/sh

#### VARs
# Destination(RSYNC)
DstDir=""			# Destination directory without ending '/'. (ex: /home/myuser)
DstSrv=""			# Hostname OR IP address of the destination. (ex: myserver.mydomain.com)
DstPort=""			# SSH Port (ex: 22)
DstUser=""			# SSH Username (ex: my_user)
DstPass=""			# SSH Password
MaxToSync=0			# Max sync to start from FilesList per cronjob ? (ex: 0 for sync all) 
# Notification (if enabled)
MailObject="MySB - Transfer to NAS completed !"

# DO NOT MODIFY
UserName="`whoami`"
ScriptName=$(basename $0)
FilesList="$HOME/scripts/files.list"
LogFile="$HOME/logs/$ScriptName.log"
LockFile="$HOME/scripts/$ScriptName.lock"
PID="$$"

# Rsync already launched: STOP
if [ -e $LockFile ] || [ ! -s $FilesList ]; then
	exit 0
fi

# Creating a lock file for do not launch two instances of in parallel
touch $LockFile

# Retrieving files to transfer
Listing() {
	# VARs
	FnLine="$1"	
	ToSynchronize="`echo $FnLine | sed 's/\ /\\\ /g;'`"
	CountSlash=`echo $FnLine |awk -F "/" '{print NF-1}'`
	(( CountSlash++ ))
	get_name="`echo $FnLine | awk '{split($0,a,"/"); print a['$CountSlash']}'`"
	get_name_escaped="`echo $get_name | sed s,/,\\\\\\\\\\/,g`"
	LastDirectory="`echo $FnLine | sed "s/$get_name_escaped//g;"`"
	CountSlash=`echo $LastDirectory |awk -F "/" '{print NF-1}'`
	LastDirectory="`echo $LastDirectory | awk '{split($0,a,"/"); print a['$CountSlash']}'`"
	CompleteDirs="`ls -A1 ~/rtorrent/complete/`"

	if [ -z "`echo $CompleteDirs | grep "$LastDirectory"`" ]; then
		FinalDst="$DstDir/"
	else
		FinalDst="$DstDir/$LastDirectory/"
	fi

	# Log
	echo "#### START ####"
	echo "From:		$FnLine"
	echo "To:		$FinalDst"
	echo "Server:		$DstUser@$DstSrv"
	echo "Date:		`/bin/date +%Y/%m/%d`"
	echo "Hour:		`/bin/date +%H:%M:%S`"
	echo "--------------------"
	echo

	# Creating sub-directory on destination
	sshpass -p $DstPass ssh -p $DstPort -o StrictHostKeyChecking=no -o UserKnownHostsFile=/dev/null $DstUser@$DstSrv "if [ ! -d $FinalDst ]; then mkdir -p $FinalDst; fi" >> $LogFileTemp

	# Start synchroniszation with RSYNC
	rsync -av --rsh="/usr/bin/sshpass -p $DstPass ssh -p$DstPort -o StrictHostKeyChecking=no -o UserKnownHostsFile=/dev/null" "$FnLine" $DstUser@$DstSrv:$FinalDst >> $LogFileTemp

	# Deleting the file from the list
	if [ $? -eq 0 ]; then
		FnLine="`echo $ToSynchronize | sed s,/,\\\\\\\\\\/,g`"
		sed -i "/$FnLine/d" $FilesList
	fi
	(( Count++ ))

	echo
	echo "--------------------"
	echo "Date:		`/bin/date +%Y/%m/%d`"
	echo "Hour:		`/bin/date +%H:%M:%S`"
	echo "##### END #####"

	# E-mail notification
	content="$LogFileTemp"
	curl --data "username=$UserName&get_custom1=$ToSynchronize&get_name=$get_name&subject=$MailObject&content=$content" http://localhost:8888/rTorrent.php
}

# Execute
Count=0
Line="`head -n 1 $FilesList | tail -n 1`"
while [ ! -z "$Line" ]; do
	LogFileTemp="/tmp/$ScriptName.log.$PID"
	Listing "$Line" >> $LogFileTemp
	# Move log file content
	if [ -f $LogFileTemp ]; then
		cat $LogFileTemp >> $LogFile
		rm -f $LogFileTemp
	fi
	Line="`head -n 1 $FilesList | tail -n 1`"
	if [ $Count -eq $MaxToSync ] && [ $MaxToSync -ne 0 ]; then
		unset Line
	fi
done

# Removing the lock file
if [ -e $LockFile ]; then
	rm -f $LockFile
fi
