# templates/tweaks/MySB_SecurityRules_Init.service.tmpl
[Unit]
Description=MySB security rules at boot
Wants=nnetworking.service mariadb.service netfilter-persistent.service
Before=remote-fs.target
After=networking.service netfilter-persistent.service ssh.service mariadb.service postfix.service

[Service]
#Type=oneshot
Type=idle
RemainAfterExit=True
ExecStart=/bin/bash <InstallDir>/bin/MySB_SecurityRules 'create' 'INIT'
ExecStop=/bin/bash <InstallDir>/bin/MySB_SecurityRules 'clean' 'STOP'

[Install]
WantedBy=multi-user.target
