# templates/systemd/MySB_SecurityRules_Init.service.tmpl
[Unit]
Description=Create all MySB security rules on boot
Before=nginx.service vsftpd.service plexmediaserver.service openvpn.service remote-fs.target
After=rsyslog.service networking.service netfilter-persistent.service ssh.service mariadb.service postfix.service
Requires=networking.service mariadb.service

[Service]
Type=oneshot
RemainAfterExit=True
ExecStart=/bin/bash /opt/MySB/bin/MySB_SecurityRules 'create' 'INIT'
ExecStop=/bin/bash /opt/MySB/bin/MySB_SecurityRules 'clean' 'STOP'

[Install]
WantedBy=multi-user.target
